# silo/ddns.timer
[Unit]
Description=Run DDNS hourly

[Timer]
OnStartupSec=15
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target